# systemd/plusclouds-agent.service
[Unit]
Description=PlusClouds Ubuntu VM Agent
Documentation=https://github.com/plusclouds/vm.agent
After=network-online.target systemd-resolved.service
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root

ExecStartPre=/bin/mkdir -p /var/log/plusclouds
ExecStartPre=/bin/chmod 0750 /var/log/plusclouds

ExecStart=/usr/local/bin/plusclouds-agent --config /etc/plusclouds/agent.yaml

Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=plusclouds-agent

# Load optional environment overrides (agent UUID, API key, log path, etc.)
EnvironmentFile=-/etc/plusclouds/environment

# Security
NoNewPrivileges=yes
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/plusclouds /etc/plusclouds
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
